# deploy/packages/linux/postinst
#!/bin/bash
/bin/systemctl start osqueryd
/bin/systemctl enable osqueryd
